# deploy/systemd/taskschmiede-proxy.service
[Unit]
Description=Taskschmiede MCP Proxy
Documentation=https://docs.taskschmiede.dev
After=network.target taskschmiede.service
Wants=taskschmiede.service

[Service]
Type=simple
User=taskschmiede
Group=taskschmiede
WorkingDirectory=/opt/taskschmiede

ExecStart=/opt/taskschmiede/bin/taskschmiede-proxy \
  --config-file /opt/taskschmiede/config/config.yaml

Restart=on-failure
RestartSec=5
TimeoutStopSec=30

EnvironmentFile=-/opt/taskschmiede/config/.env

StandardOutput=journal
StandardError=journal
SyslogIdentifier=taskschmiede-proxy

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/var/log/taskschmiede
ReadOnlyPaths=/opt/taskschmiede/bin /opt/taskschmiede/config

[Install]
WantedBy=multi-user.target
