# deploy/systemd/taskschmiede-portal.service
[Unit]
Description=Taskschmiede Portal (Web UI)
Documentation=https://docs.taskschmiede.dev
After=network.target taskschmiede.service
Wants=taskschmiede.service

[Service]
Type=simple
User=taskschmiede
Group=taskschmiede
WorkingDirectory=/opt/taskschmiede

EnvironmentFile=-/opt/taskschmiede/config/.env

ExecStart=/opt/taskschmiede/bin/taskschmiede-portal \
  --config /opt/taskschmiede/config/config.yaml \
  --secure

Restart=on-failure
RestartSec=5
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=taskschmiede-portal

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadOnlyPaths=/opt/taskschmiede/bin /opt/taskschmiede/config

[Install]
WantedBy=multi-user.target
